# deploy/systemd/smartai-backup.service
[Unit]
Description=Backup do banco do SmartAI
After=network.target

[Service]
Type=oneshot
User=smartai
Group=smartai
WorkingDirectory=/opt/smartai/app
EnvironmentFile=/opt/smartai/app/.env.production
ExecStart=/opt/smartai/app/deploy/scripts/backup_postgres.sh
